# health-concierge/scripts/health-concierge.service
[Unit]
Description=Health Concierge Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/health-product-analysis/health-concierge
ExecStart=/root/health-product-analysis/health-concierge/.venv/bin/python -m src.bot
Restart=on-failure
RestartSec=10
StandardOutput=append:/root/health-product-analysis/health-concierge/logs/bot.log
StandardError=append:/root/health-product-analysis/health-concierge/logs/bot.log

[Install]
WantedBy=multi-user.target
